# roles/systemd-networkd/files/rootfs/usr/lib/systemd/system/ndp-fix@.timer
[Unit]
Description=Runs the ndp-fix script on IP %i

[Timer]
OnUnitActiveSec=30
RandomizedDelaySec=10