# roles/rpki-client/files/rootfs/lib/systemd/system/rpki-client.timer
[Unit]
Description=Periodically run rpki-client

[Timer]
OnCalendar=hourly
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target